# deploy/systemd/gracchus-projects-contracts.timer
[Unit]
Description=Gracchus — projects & contracts refresh (daily 04:15 UTC)

[Timer]
OnCalendar=*-*-* 04:15:00
RandomizedDelaySec=300
Persistent=true
Unit=gracchus-job@projects-contracts.service

[Install]
WantedBy=timers.target
